# systemd/ai-prorok-generate.timer
[Unit]
Description=AI Prorok - Generation Timer
Requires=ai-prorok-generate.service

[Timer]
# Проверка каждую минуту, запуск если время совпадает
OnCalendar=*-*-* *:*:00
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
