# systemd/ai-prorok-publish.service
[Unit]
Description=AI Prorok - Publication Service
After=network.target

[Service]
Type=oneshot
User=YOUR_USERNAME
Group=YOUR_GROUP
WorkingDirectory=/path/to/ai_prorok
Environment="PATH=/path/to/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Сначала проверяем, пора ли публиковать
ExecStartPre=/path/to/venv/bin/python3 /path/to/ai_prorok/scheduler.py check-publish

# Если проверка прошла (exit 0), публикуем
ExecStart=/path/to/venv/bin/python3 /path/to/ai_prorok/ai_prorok_refactored.py publish

# Настройки безопасности
NoNewPrivileges=true
PrivateTmp=true

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-prorok-publish

[Install]
WantedBy=multi-user.target
