# systemd/ai-prorok-publish.timer
[Unit]
Description=AI Prorok - Publication Timer
Requires=ai-prorok-publish.service

[Timer]
# Проверка каждую минуту, запуск если время совпадает
OnCalendar=*-*-* *:*:00
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
